# device/copy/etc/supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/run/supervisor.pid
logfile=/var/log/supervisord.log
loglevel=error

[unix_http_server]
file=/run/supervisord.sock
username=sv
password=sv

[supervisorctl]
serverurl=unix:///run/supervisord.sock
username=sv
password=sv

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
startretries=0
autorestart=true
stdout_logfile=/var/log/php-fpm.out.log
stderr_logfile=/var/log/php-fpm.err.log
priority=100

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
startretries=0
autorestart=true
stdout_logfile=/var/log/sshd.out.log
stderr_logfile=/var/log/sshd.err.log
priority=101

[program:phd]
user=phab
command=sh -c "/code/phabricator/scripts/daemon/phd-daemon < /usr/local/etc/phd-config.json"
autostart=true
startretries=0
autorestart=false
stdout_logfile=/var/log/phd-daemon.out.log
stderr_logfile=/var/log/phd-daemon.err.log
priority=102
